# invenio-supervisor.conf
[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── InvenioRDM UI (uWSGI :5000) ───────────────────────────────────────────────
[program:invenio-ui]
command=/home/%(ENV_USER)s/invenio-venv/bin/uwsgi
        /home/%(ENV_USER)s/invenio-instance/uwsgi_ui.ini
environment=INVENIO_INSTANCE_PATH="/home/%(ENV_USER)s/invenio-instance"
directory=/home/%(ENV_USER)s/invenio-rdm/src
autostart=true
autorestart=true
stdout_logfile=/tmp/invenio-ui.log
stderr_logfile=/tmp/invenio-ui.log

; ── InvenioRDM REST API (uWSGI :5001) ────────────────────────────────────────
[program:invenio-api]
command=/home/%(ENV_USER)s/invenio-venv/bin/uwsgi
        /home/%(ENV_USER)s/invenio-instance/uwsgi_rest.ini
environment=INVENIO_INSTANCE_PATH="/home/%(ENV_USER)s/invenio-instance"
directory=/home/%(ENV_USER)s/invenio-rdm/src
autostart=true
autorestart=true
stdout_logfile=/tmp/invenio-api.log
stderr_logfile=/tmp/invenio-api.log

; ── Celery worker ─────────────────────────────────────────────────────────────
[program:celery]
command=/home/%(ENV_USER)s/invenio-venv/bin/celery
        -A invenio_app.celery worker --loglevel=INFO --concurrency=2
environment=INVENIO_INSTANCE_PATH="/home/%(ENV_USER)s/invenio-instance"
directory=/home/%(ENV_USER)s/invenio-rdm/src
autostart=true
autorestart=true
stdout_logfile=/tmp/celery.log
stderr_logfile=/tmp/celery.log

; ── MinIO object storage ──────────────────────────────────────────────────────
[program:minio]
command=minio server /home/%(ENV_USER)s/minio/data --console-address ":9001"
environment=MINIO_ROOT_USER="minio",MINIO_ROOT_PASSWORD="minio123456"
autostart=true
autorestart=true
stdout_logfile=/tmp/minio.log
stderr_logfile=/tmp/minio.log
